# ops/systemd/tracey-customer-agent.service
[Unit]
Description=Tracey Customer Voice Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/tracey-agent/current/livekit-agent
EnvironmentFile=/etc/tracey-voice/common.env
EnvironmentFile=/etc/tracey-voice/deploy.env
Environment=VOICE_WORKER_ROLE=tracey-customer-agent
Environment=VOICE_WORKER_SURFACES=normal
ExecStart=/usr/bin/npm run start:customer
Restart=always
RestartSec=5
SyslogIdentifier=tracey-customer-agent
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
